# services/heartbeat.service
[Unit]
Description=P2BP Heartbeat Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/p2bp/camera/scripts/heartbeat.py
WorkingDirectory=/opt/p2bp/camera
Restart=always
RestartSec=5

# Ensure runtime signal dir exists
RuntimeDirectory=p2bp
RuntimeDirectoryMode=0755

# Environment
EnvironmentFile=/opt/p2bp/camera/config/agent.env

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
