# services/homography.service
[Unit]
Description=P2BP Homography Agent
After=network-online.target camera-scanner.service
Wants=network-online.target
Requires=camera-scanner.service

[Service]
Type=simple
WorkingDirectory=/opt/p2bp/camera
ExecStart=/usr/bin/python3 -s -m scripts.homography

Restart=always
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

ReadWritePaths=/opt/p2bp/camera /run/p2bp

RuntimeDirectory=p2bp
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
