# services/scan-orchestrator.service
[Unit]
Description=P2BP Lidar Scan Orchestrator (poll Web, upload .xyz)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -s -m scripts.scan_orchestrator
WorkingDirectory=/opt/p2bp/camera
Restart=always
RestartSec=10

LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

EnvironmentFile=/opt/p2bp/camera/config/agent.env
# LidarScanV1 uses timed pauses instead of input() when non-interactive
Environment=P2BP_LIDAR_SCAN_NONINTERACTIVE=1

# When the lidar runs on a Raspberry Pi and this host (Jetson) only polls Web and uploads:
#   Environment=P2BP_SCAN_CMD=/opt/p2bp/camera/scripts/run_lidar_on_pi.sh
# (run_lidar_on_pi.sh defaults to pi@192.168.28.2; set P2BP_PI_SSH in agent.env to override.)
# Optional: Environment=P2BP_PI_REMOTE_SCRIPT=/opt/p2bp/lidar/scripts/LidarScanV1.py
# Optional: Environment=P2BP_PI_REMOTE_XYZ=/tmp/p2bp_lidar_scan.xyz

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
