# services/tracker.service
[Unit]
Description=P2BP Tracking Agent
After=network-online.target
Wants=network-online.target
ConditionPathExists=/run/p2bp/tracking.enabled

[Service]
Type=simple
WorkingDirectory=/opt/p2bp/camera
ExecStart=/usr/bin/python3 -s -m scripts.tracker

Environment=PYTHONUNBUFFERED=1
Environment=HOME=/opt/p2bp/camera
Environment=XDG_CACHE_HOME=/opt/p2bp/camera/.cache
Environment=MPLCONFIGDIR=/opt/p2bp/camera/.cache/matplotlib
Environment=YOLO_CONFIG_DIR=/opt/p2bp/camera/.config/Ultralytics

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/p2bp/camera

Restart=no
KillSignal=SIGTERM
TimeoutStopSec=10

RuntimeDirectory=p2bp
RuntimeDirectoryMode=0755

EnvironmentFile=/opt/p2bp/camera/config/agent.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
