# services/tracker.path
[Unit]
Description=Enable tracking script when flag exists

[Path]
PathExists=/run/p2bp/tracking.enabled
PathChanged=/run/p2bp/tracking.enabled
Unit=tracker-toggle.service

[Install]
WantedBy=multi-user.target
